# resources/arrakis-cdpserver.service
[Unit]
Description=Chrome DevTools Protocol Server for Arrakis
After=arrakis-guestinit.service

[Service]
Type=simple
User=elara
WorkingDirectory=/home/elara
ExecStart=/usr/local/bin/arrakis-cdpserver
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
