# resources/arrakis-chrome-forwarder.service
[Unit]
Description=Chrome CDP Port Forwarder
After=arrakis-chrome.service
Requires=arrakis-chrome.service

[Service]
Type=simple
User=elara
Group=elara
ExecStart=/usr/bin/socat TCP-LISTEN:9223,bind=0.0.0.0,fork TCP:127.0.0.1:9222
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
