# resources/arrakis-novncserver.service
[Unit]
Description=NoVNC Server for Arrakis
After=arrakis-vncserver.service
Requires=arrakis-vncserver.service

[Service]
Type=simple
User=elara
WorkingDirectory=/opt/novnc
ExecStart=/usr/local/bin/start-novnc
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
