# resources/arrakis-novncserver.service
[Unit]
Description=NoVNC Server for Arrakis
After=arrakis-vncserver.service
Requires=arrakis-vncserver.service

[Service]
Type=simple
User=elara
WorkingDirectory=/opt/novnc
ExecStart=/opt/novnc/utils/novnc_proxy --listen 6080 --vnc localhost:5901
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
